# ansible/roles/wheresalice.sshpoints/templates/ssh.service.j2
[Unit]
Description=sshpoints ssh service
Requires=network.target
After=syslog.target network.target
[Service]
Type=simple
ExecStart=/srv/sshpoints/sshpoints s
[Install]
WantedBy=multi-user.target
